# web_svc/atrax_keeper/deployment/atrax_keeper.conf
start on runlevel [2345]
stop on runlevel [!2345]

console log

exec /usr/local/etc/atrax/{{ module }}/{{ module }}/run.sh {{ crawl_job }}
